# overlay/control/postinst
#! /bin/bash

# Q1Libre: fix sudoers permissions immediately (dpkg extracts with Windows perms)
chmod 0440 /etc/sudoers.d/q1libre 2>/dev/null || true

python3 /root/scripts/merge.py

# Q1Libre: Re-enable update_manager after merge.py may have reset it
sed -i 's/^# \[update_manager\]/[update_manager]/' /home/mks/klipper_config/moonraker.conf
sed -i 's/^# channel: dev/channel: dev/' /home/mks/klipper_config/moonraker.conf
sed -i 's/^# refresh_interval: 168/refresh_interval: 168/' /home/mks/klipper_config/moonraker.conf
sed -i 's/^# enable_system_updates: False/enable_system_updates: False/' /home/mks/klipper_config/moonraker.conf
sed -i 's/^# \[update_manager fluidd\]/[update_manager fluidd]/' /home/mks/klipper_config/moonraker.conf
sed -i '/^\[update_manager fluidd\]/,/^$/{s/^# type: web/type: web/; s/^# channel: stable/channel: stable/; s|^# repo: fluidd-core/fluidd|repo: fluidd-core/fluidd|; s|^# path: ~/fluidd|path: ~/fluidd|}' /home/mks/klipper_config/moonraker.conf

if [ ! -f /etc/bash.bashrc ]; then
    echo "Missing bash.bashrc file, recovering..."
    cp -fp /root/etc/bash.bashrc /etc/bash.bashrc
fi

# Read connection method (0=LAN, 1=Internet)
method_value=0
while IFS='=' read -r key value; do
    key=$(echo $key | xargs)
    value=$(echo $value | xargs)
    if [ "$key" = "method" ]; then
        method_value=$value
        break
    fi
done < <(sed -n '/\[app_connection\]/,/\[/p' /home/mks/klipper_config/config.mksini)

# frpc install if missing
if ! systemctl list-unit-files | grep -qw '^frpc.service'; then
    chmod +x /root/frp/setup_frpc.sh /root/frp/frpc
    bash /root/frp/setup_frpc.sh
fi

# Comment out Armbian source (no longer supported for Buster)
sudo sed -i '1s/^/#/' /etc/apt/sources.list.d/armbian.list

# Q1Libre: Remove hardcoded Chinese DNS, use DHCP-provided DNS with public fallbacks
echo "Q1Libre: fixing DNS configuration..."
cat > /etc/resolv.conf <<'DNSEOF'
nameserver 1.1.1.1
nameserver 8.8.8.8
DNSEOF
if [ -f /etc/NetworkManager/NetworkManager.conf ]; then
    sed -i 's/^rc-manager=file/rc-manager=resolvconf/' /etc/NetworkManager/NetworkManager.conf
    systemctl restart NetworkManager 2>/dev/null || true
fi

rm -f /etc/profile.d/time_refresh.sh
mv -f /root/auto_update/auto_update_logrotate /etc/logrotate.d/auto_update_logrotate

# Sudoers for mks user
SUDO_FILE="/etc/sudoers.d/mks"
sudo systemctl reload nginx
echo 'mks ALL=(ALL) ALL, !/usr/sbin/visudo, !/bin/nano /etc/sudoers, !/bin/nano /root/*, !/bin/ls /root, !/bin/cat /etc/sudoers, !/bin/systemctl, !/bin/chmod, NOPASSWD: /usr/bin/psd-overlay-helper, /usr/bin/timedatectl, /usr/bin/apt' | sudo tee $SUDO_FILE > /dev/null
chmod 0440 $SUDO_FILE

# Adaptive mesh config
if [ ! -f "/home/mks/klipper_config/Adaptive_Mesh.cfg" ]; then
    if [ -f "/root/Adaptive_Mesh.cfg" ]; then
        cp /root/Adaptive_Mesh.cfg /home/mks/klipper_config/ -f
        chmod 644 /home/mks/klipper_config/Adaptive_Mesh.cfg
    fi
fi

cp -f /root/10-armbian-header /etc/update-motd.d/10-armbian-header

################本地更新python库
local_package_dir="/home/mks/python_package"
installed_packages=$(pip3 list)

install_python_package_if_not_exists() {
    package_name=$1
    package_file=$2
    if ! echo "$installed_packages" | grep -Fq $package_name; then
        echo "Installing $package_name..."
        pip3 install --no-index --find-links "$local_package_dir" "$local_package_dir/$package_file"
    fi
}

need_update=0
for package in "qrcode" "cryptography" "toml"; do
    if ! echo "$installed_packages" | grep -Fq $package; then
        need_update=1
        break
    fi
done

if [ $need_update -eq 1 ]; then
    python3 -m pip install --no-index --find-links "pip" "/home/mks/python_package/pip-23.3.2-py3-none-any.whl"
    python3 -m pip install --no-index --find-links "setuptools" "/home/mks/python_package/setuptools-68.0.0-py3-none-any.whl"
    python3 -m pip install --no-index --find-links "wheel" "/home/mks/python_package/wheel-0.42.0-py3-none-any.whl"
fi

install_python_package_if_not_exists "qrcode" "qrcode-7.4.2-py3-none-any.whl"
install_python_package_if_not_exists "cryptography" "cryptography-41.0.7-cp37-abi3-manylinux_2_17_aarch64.manylinux2014_aarch64.whl"
install_python_package_if_not_exists "toml" "toml-0.10.2-py2.py3-none-any.whl"
rm -rf /home/mks/python_package
################本地更新python库结束

# Copy hardware config files
cp -fp /root/etc/rk3328-roc-cc.dtb /boot/dtb/rockchip/rk3328-roc-cc.dtb
cp -fp /root/etc/rk3328-roc-cc.dts /boot/dtb/rockchip/rk3328-roc-cc.dts
cp -fp /root/etc/60-persistent-v4l.rules /lib/udev/rules.d/60-persistent-v4l.rules

# Fix permissions (Q1Libre: proper 755/644 instead of stock 777)
find /home/mks/klipper_config -type d -exec chmod 755 {} \;
find /home/mks/klipper_config -type f -exec chmod 644 {} \;
chmod 700 /root/xindi/build/xindi
chown -R mks:mks /home/mks/klipper/
chown -R mks:mks /home/mks/klipper_config/
chmod 700 /usr/bin/makerbase-automount
chmod 0440 /etc/sudoers.d/q1libre 2>/dev/null || true

# Q1Libre: Fix stock 777/666 permissions on safe directories
echo "Q1Libre: fixing stock file permissions..."
find /home/mks/klipper -type d -exec chmod 755 {} \;
find /home/mks/klipper -type f -exec chmod 644 {} \;
find /home/mks/klipper/scripts -type f -name "*.sh" -exec chmod 755 {} \;
find /home/mks/moonraker -type d -exec chmod 755 {} \;
find /home/mks/moonraker -type f -exec chmod 644 {} \;
find /home/mks/moonraker/scripts -type f -name "*.sh" -exec chmod 755 {} \;
find /home/mks/fluidd -type d -exec chmod 755 {} \; 2>/dev/null
find /home/mks/fluidd -type f -exec chmod 644 {} \; 2>/dev/null
find /root/scripts -type f -exec chmod 644 {} \; 2>/dev/null
chmod 755 /root/scripts/*.sh 2>/dev/null
find /root/etc -type f -exec chmod 644 {} \; 2>/dev/null
chmod 644 /root/xindi/version /root/xindi/mks-id.sh 2>/dev/null

# Q1Libre: install version marker
if [ -f /root/q1libre_version.txt ]; then
    echo "Q1Libre firmware installed" >> /var/log/q1libre_install.log
    echo "  Version: $(cat /root/q1libre_version.txt)" >> /var/log/q1libre_install.log
    echo "  Date: $(date)" >> /var/log/q1libre_install.log
fi

# ═══════════════════════════════════════════════════════════════
# Q1Libre: Klipper upgrade
# ═══════════════════════════════════════════════════════════════
echo "Q1Libre: upgrading Klipper..."
systemctl stop klipper.service 2>/dev/null || true

# Back up existing klipper
if [ -d /home/mks/klipper ]; then
    rm -rf /home/mks/klipper.bak
    cp -a /home/mks/klipper /home/mks/klipper.bak \
        || echo "Q1Libre: WARNING - klipper backup failed"
fi

chown -R mks:mks /home/mks/klipper

# OPTIMIZATION: Skip klippy-env rebuild if already Python 3 with correct packages
NEED_VENV_REBUILD=0
if [ ! -x /home/mks/klippy-env/bin/python3 ]; then
    NEED_VENV_REBUILD=1
elif ! /home/mks/klippy-env/bin/python3 -c "import numpy, greenlet, jinja2" 2>/dev/null; then
    NEED_VENV_REBUILD=1
fi

if [ $NEED_VENV_REBUILD -eq 1 ]; then
    echo "Q1Libre: rebuilding klippy-env with Python 3..."
    if [ -d /home/mks/klippy-env ]; then
        cp -a /home/mks/klippy-env /home/mks/klippy-env.bak
        rm -rf /home/mks/klippy-env
    fi
    /usr/bin/virtualenv -p python3 /home/mks/klippy-env 2>/dev/null
    if [ -x /home/mks/klippy-env/bin/pip ]; then
        if [ -d /root/klippy-wheels ] && ls /root/klippy-wheels/*.whl >/dev/null 2>&1; then
            echo "Q1Libre: installing klippy packages from bundled wheels..."
            if ls /root/klippy-wheels/pip-*.whl >/dev/null 2>&1; then
                /home/mks/klippy-env/bin/pip install --no-index /root/klippy-wheels/pip-*.whl 2>&1 || true
            fi
            if ls /root/klippy-wheels/setuptools-*.whl >/dev/null 2>&1; then
                /home/mks/klippy-env/bin/pip install --no-index /root/klippy-wheels/setuptools-*.whl 2>&1 || true
            fi
            /home/mks/klippy-env/bin/pip install --no-index --find-links /root/klippy-wheels \
                /root/klippy-wheels/greenlet-*.whl \
                /root/klippy-wheels/cffi-*.whl \
                /root/klippy-wheels/pycparser-*.whl \
                /root/klippy-wheels/Jinja2-*.whl \
                /root/klippy-wheels/MarkupSafe-*.whl \
                /root/klippy-wheels/pyserial-*.whl \
                /root/klippy-wheels/wrapt-*.whl \
                /root/klippy-wheels/aenum-*.whl \
                /root/klippy-wheels/numpy-*.whl \
                /root/klippy-wheels/scipy-*.whl \
                2>&1 || echo "Q1Libre: WARNING - some wheel installs had errors"
            if ls /root/klippy-wheels/python-can-*.tar.gz >/dev/null 2>&1; then
                /home/mks/klippy-env/bin/pip install --no-index --find-links /root/klippy-wheels \
                    /root/klippy-wheels/python-can-*.tar.gz 2>&1 || true
            fi
            rm -rf /root/klippy-wheels 2>/dev/null || true
        else
            /home/mks/klippy-env/bin/pip install -r /home/mks/klipper/scripts/klippy-requirements.txt \
                --quiet 2>/dev/null \
                || echo "Q1Libre: WARNING - klippy pip install had errors"
        fi
        chown -R mks:mks /home/mks/klippy-env
        echo "Q1Libre: klippy-env rebuilt"
    else
        echo "Q1Libre: WARNING - klippy-env rebuild failed, restoring backup"
        rm -rf /home/mks/klippy-env
        cp -a /home/mks/klippy-env.bak /home/mks/klippy-env
    fi
else
    echo "Q1Libre: klippy-env already correct, skipping rebuild"
fi

# OPTIMIZATION: Use pre-compiled c_helper.so if bundled, skip compilation
if [ -f /home/mks/klipper/klippy/chelper/c_helper.so ]; then
    echo "Q1Libre: using bundled c_helper.so"
    chown mks:mks /home/mks/klipper/klippy/chelper/c_helper.so
    chmod 755 /home/mks/klipper/klippy/chelper/c_helper.so
else
    echo "Q1Libre: compiling klipper C helper..."
    cd /home/mks/klipper/klippy/chelper
    /home/mks/klippy-env/bin/python -c "
import sys; sys.path.insert(0, '..')
import chelper; chelper.get_ffi()
" 2>&1 && echo "Q1Libre: c_helper compiled" || echo "Q1Libre: WARNING - c_helper compilation failed"
    chown mks:mks /home/mks/klipper/klippy/chelper/c_helper.so 2>/dev/null || true
    cd /
fi

# Git repo setup for klipper
KLIPPER_VERSION="v0.13.0-q1pro"
KLIPPER_SHA="7740b479e94084866a74a5b0c4a59d566aa90d08"
cd /home/mks/klipper
rm -rf .git
git init 2>/dev/null
git config user.email "q1libre@local"
git config user.name "Q1Libre"
git remote add origin https://github.com/danthi123/klipper.git 2>/dev/null || true
git add -A 2>/dev/null
git commit -m "Q1Libre: vendored klipper $KLIPPER_VERSION" --quiet 2>/dev/null
git tag "$KLIPPER_VERSION" 2>/dev/null

# OPTIMIZATION: Run git fetch in background (non-blocking)
(git fetch origin --tags 2>/dev/null && \
 git cat-file -e $KLIPPER_SHA 2>/dev/null && \
 git reset --hard $KLIPPER_SHA 2>/dev/null && \
 git branch --set-upstream-to=origin/q1-pro master 2>/dev/null && \
 chown -R mks:mks /home/mks/klipper/.git && \
 echo "Q1Libre: klipper git fetch complete") &
KLIPPER_GIT_PID=$!

chown -R mks:mks /home/mks/klipper/.git
cd /

# Link timelapse macros
if [ -f /home/mks/moonraker-timelapse/klipper_macro/timelapse.cfg ]; then
    ln -sf /home/mks/moonraker-timelapse/klipper_macro/timelapse.cfg /home/mks/klipper_config/timelapse.cfg
    chown -h mks:mks /home/mks/klipper_config/timelapse.cfg
fi

echo "Q1Libre: Klipper upgrade complete"

# ═══════════════════════════════════════════════════════════════
# Q1Libre: Moonraker upgrade
# ═══════════════════════════════════════════════════════════════
echo "Q1Libre: upgrading Moonraker..."
systemctl stop moonraker.service 2>/dev/null || true

if [ -d /home/mks/moonraker ]; then
    rm -rf /home/mks/moonraker.bak
    cp -a /home/mks/moonraker /home/mks/moonraker.bak \
        || echo "Q1Libre: WARNING - moonraker backup failed"
fi

chown -R mks:mks /home/mks/moonraker

# Git repo setup for moonraker
MOONRAKER_VERSION="v0.10.0-19"
VENDORED_SHA="1ed102edfb34906115bfeb02712d4c839b93a7a9"
cd /home/mks/moonraker
rm -rf .git
git init 2>/dev/null
git config user.email "q1libre@local"
git config user.name "Q1Libre"
git remote add origin https://github.com/Arksine/moonraker.git 2>/dev/null || true
git add -A 2>/dev/null
git commit -m "Q1Libre: vendored moonraker $MOONRAKER_VERSION" --quiet 2>/dev/null
git tag "$MOONRAKER_VERSION" 2>/dev/null
mkdir -p .git/info
echo "/thirdparty" >> .git/info/exclude

# OPTIMIZATION: Run git fetch in background (non-blocking)
(git fetch origin --tags 2>/dev/null && \
 git cat-file -e $VENDORED_SHA 2>/dev/null && \
 git reset --hard $VENDORED_SHA 2>/dev/null && \
 git branch --set-upstream-to=origin/master master 2>/dev/null && \
 echo "/moonraker/components/timelapse.py" >> /home/mks/moonraker/.git/info/exclude 2>/dev/null && \
 ln -sf /home/mks/moonraker-timelapse/component/timelapse.py /home/mks/moonraker/moonraker/components/timelapse.py 2>/dev/null && \
 chown -h mks:mks /home/mks/moonraker/moonraker/components/timelapse.py 2>/dev/null && \
 chown -R mks:mks /home/mks/moonraker/.git && \
 echo "Q1Libre: moonraker git fetch complete") &
MOONRAKER_GIT_PID=$!

# Clear moonraker's cached version data so it reads fresh git state
rm -f /home/mks/.moonraker_database/data.mdb /home/mks/.moonraker_database/lock.mdb 2>/dev/null

chown -R mks:mks /home/mks/moonraker/.git
cd /

# Create thirdparty symlink
ln -sf /home/mks/moonraker/moonraker/thirdparty /home/mks/moonraker/thirdparty
chown -h mks:mks /home/mks/moonraker/thirdparty

# Link timelapse plugin
if [ -f /home/mks/moonraker-timelapse/component/timelapse.py ]; then
    ln -sf /home/mks/moonraker-timelapse/component/timelapse.py /home/mks/moonraker/moonraker/components/timelapse.py
    chown -h mks:mks /home/mks/moonraker/moonraker/components/timelapse.py
fi

echo "Q1Libre: Moonraker upgrade complete"

# printer_data symlinks
mkdir -p /home/mks/printer_data
[ -L /home/mks/printer_data/config ] || (rm -rf /home/mks/printer_data/config && ln -sf /home/mks/klipper_config /home/mks/printer_data/config)
[ -L /home/mks/printer_data/gcodes ] || (rm -rf /home/mks/printer_data/gcodes && ln -sf /home/mks/gcode_files /home/mks/printer_data/gcodes)
[ -L /home/mks/printer_data/logs ] || (rm -rf /home/mks/printer_data/logs && ln -sf /home/mks/klipper_logs /home/mks/printer_data/logs)
chown -h mks:mks /home/mks/printer_data /home/mks/printer_data/config /home/mks/printer_data/gcodes /home/mks/printer_data/logs

# Enable mDNS/Avahi
if command -v avahi-daemon >/dev/null 2>&1; then
    systemctl enable avahi-daemon 2>/dev/null || true
    systemctl restart avahi-daemon 2>/dev/null || true
fi

# ═══════════════════════════════════════════════════════════════
# Restart services — get the printer online ASAP
# ═══════════════════════════════════════════════════════════════
systemctl daemon-reload
systemctl restart klipper.service
systemctl restart moonraker.service
if [ "$method_value" = "1" ]; then
    systemctl enable frpc.service
    systemctl start frpc.service
fi

# OPTIMIZATION: Reduced from 12s to 5s — moonraker starts in ~3s
sleep 5
curl -s -X POST http://127.0.0.1:7125/printer/firmware_restart >/dev/null 2>&1

echo "Q1Libre: services restarted"

# ═══════════════════════════════════════════════════════════════
# Q1Libre: printer.cfg fixes (safe idempotent patches)
# ═══════════════════════════════════════════════════════════════
PCFG=/home/mks/klipper_config/printer.cfg
if [ -f "$PCFG" ]; then
    # Enable filament runout sensor (Qidi ships it disabled)
    sed -i '/^\[hall_filament_width_sensor\]/,/^\[/{s/^enable: false/enable: true/}' "$PCFG"
    # Loosen probe tolerance from 0.05 to 0.10 (Klipper default, avoids excessive retries)
    sed -i 's/^samples_tolerance: 0\.05$/samples_tolerance: 0.10/' "$PCFG"
    echo "Q1Libre: printer.cfg patched (filament sensor, probe tolerance)"
fi

# ═══════════════════════════════════════════════════════════════
# Post-restart tasks (run AFTER services are up so user gets a working printer fast)
# ═══════════════════════════════════════════════════════════════

# Wait for background git fetches to complete (max 30s timeout)
if [ -n "$KLIPPER_GIT_PID" ]; then
    wait $KLIPPER_GIT_PID 2>/dev/null || true
fi
if [ -n "$MOONRAKER_GIT_PID" ]; then
    wait $MOONRAKER_GIT_PID 2>/dev/null || true
fi

# OPTIMIZATION: Debloat runs AFTER services are restarted (non-blocking for user)
if [ -x /root/scripts/debloat.sh ]; then
    echo "Q1Libre: running system debloat in background..."
    bash /root/scripts/debloat.sh &
fi

# Security updates (online only, runs in background)
if [ "$method_value" = "1" ]; then
    (apt-get update -qq 2>/dev/null && \
     apt-get upgrade -y -qq 2>/dev/null && \
     echo "Q1Libre: security updates applied") &
fi

sync
echo "Q1Libre: install complete"
