# overlay/control/postinst
#! /bin/bash

# Q1Libre: fix sudoers permissions immediately (dpkg extracts with Windows perms)
chmod 0440 /etc/sudoers.d/q1libre 2>/dev/null || true

python3 /root/scripts/merge.py

# Q1Libre: Re-enable update_manager after merge.py may have reset it
sed -i 's/^# \[update_manager\]/[update_manager]/' /home/mks/klipper_config/moonraker.conf
sed -i 's/^# channel: dev/channel: dev/' /home/mks/klipper_config/moonraker.conf
sed -i 's/^# refresh_interval: 168/refresh_interval: 168/' /home/mks/klipper_config/moonraker.conf
sed -i 's/^# enable_system_updates: False/enable_system_updates: False/' /home/mks/klipper_config/moonraker.conf
sed -i 's/^# \[update_manager fluidd\]/[update_manager fluidd]/' /home/mks/klipper_config/moonraker.conf
sed -i '/^\[update_manager fluidd\]/,/^$/{s/^# type: web/type: web/; s/^# channel: stable/channel: stable/; s|^# repo: fluidd-core/fluidd|repo: fluidd-core/fluidd|; s|^# path: ~/fluidd|path: ~/fluidd|}' /home/mks/klipper_config/moonraker.conf

# 检查文件
if [ ! -f /etc/bash.bashrc ]; then
    echo "Misiing bash.bashrc file, recovering..."
    cp -fp /root/etc/bash.bashrc /etc/bash.bashrc
fi

# 初始化变量
method_value=0

# 从config.mksini文件中读取method值
while IFS='=' read -r key value; do
    # 去除前后空格
    key=$(echo $key | xargs)
    value=$(echo $value | xargs)

    if [ "$key" = "method" ]; then
        method_value=$value
        break
    fi
done < <(sed -n '/\[app_connection\]/,/\[/p' /home/mks/klipper_config/config.mksini)

# frpc如果不存在则安装
if ! systemctl list-unit-files | grep -qw '^frpc.service'; then
    chmod +x /root/frp/setup_frpc.sh /root/frp/frpc
    bash /root/frp/setup_frpc.sh
fi

# 由于armbian官方的更新 buster似乎不再被支持，因此先注释掉armbian源
sudo sed -i '1s/^/#/' /etc/apt/sources.list.d/armbian.list

#删除旧的时间更新脚本
rm -f /etc/profile.d/time_refresh.sh

#添加日志轮换
mv -f /root/auto_update/auto_update_logrotate /etc/logrotate.d/auto_update_logrotate

# 修改sudo权限
SUDO_FILE="/etc/sudoers.d/mks"

# 重启nginx
sudo systemctl reload nginx


# 使用 echo 和 sudo tee 命令来写入文件内容
echo 'mks ALL=(ALL) ALL, !/usr/sbin/visudo, !/bin/nano /etc/sudoers, !/bin/nano /root/*, !/bin/ls /root, !/bin/cat /etc/sudoers, !/bin/systemctl, !/bin/chmod, NOPASSWD: /usr/bin/psd-overlay-helper, /usr/bin/timedatectl, /usr/bin/apt' | sudo tee $SUDO_FILE > /dev/null
# 设置正确的文件权限
chmod 0440 $SUDO_FILE

if [ ! -f "/home/mks/klipper_config/Adaptive_Mesh.cfg" ]; then
    if [ -f "/root/Adaptive_Mesh.cfg" ]; then
        cp /root/Adaptive_Mesh.cfg /home/mks/klipper_config/ -f
        chmod 644 /home/mks/klipper_config/Adaptive_Mesh.cfg
    fi
fi

#更改ssh登录界面的样式
cp -f /root/10-armbian-header /etc/update-motd.d/10-armbian-header

##仅在互联网情况下更新
#if [ "$method_value" = "1" ]; then
################更新qrcode用到的python包
#    INSTALLED_PACKAGES=$(pip3 list)
#    UPDATE_PIP=0
#    check_package() {
#        PACKAGE=$1
#        echo "Checking $PACKAGE..."
#        if echo "$INSTALLED_PACKAGES" | grep -F $PACKAGE > /dev/null; then
#            echo "$PACKAGE is already installed."
#        else
#            echo "$PACKAGE is not installed."
#            UPDATE_PIP=1  # 标记需要更新pip
#        fi
#    }
#
#    # 检查每个包
#    check_package toml
#    check_package qrcode
#    check_package cryptography
#
#    # 如果任何包没有安装，则更新pip3
#    if [ $UPDATE_PIP -eq 1 ]; then
#        echo "Updating pip3..."
#        pip3 install --upgrade pip
#    fi
#
#    install_package() {
#        PACKAGE=$1
#        echo "Installing $PACKAGE..."
#        pip3 install $PACKAGE
#    }
#
#    # 检查并安装每个包，这次不检查pip是否需要更新
#    for PACKAGE in toml qrcode cryptography; do
#        if echo "$INSTALLED_PACKAGES" | grep -F $PACKAGE > /dev/null; then
#            echo "$PACKAGE is already installed."
#        else
#            install_package $PACKAGE
#        fi
#    done
#
#    echo "All necessary packages are installed."
################更新qrcode结束
#fi

################本地更新python库
# 指定本地包的存储目录
local_package_dir="/home/mks/python_package"

# 获取已安装的Python包列表
installed_packages=$(pip3 list)

# 定义一个函数来安装Python包，如果它们尚未安装
install_python_package_if_not_exists() {
    package_name=$1
    package_file=$2

    # 检查包是否已安装
    if ! echo "$installed_packages" | grep -Fq $package_name; then
        # 如果未安装，则安装包
        echo "Installing $package_name..."
        pip3 install --no-index --find-links "$local_package_dir" "$local_package_dir/$package_file"
    else
        echo "$package_name is already installed."
    fi
}

# 检查是否需要更新 pip, setuptools, 和 wheel
need_update=0
for package in "qrcode" "cryptography" "toml"; do
    if ! echo "$installed_packages" | grep -Fq $package; then
        need_update=1
        break
    fi
done

# 如果任何一个未安装，先更新 pip, setuptools, 和 wheel
if [ $need_update -eq 1 ]; then
    python3 -m pip install --no-index --find-links "pip" "/home/mks/python_package/pip-23.3.2-py3-none-any.whl"
    python3 -m pip install --no-index --find-links "setuptools" "/home/mks/python_package/setuptools-68.0.0-py3-none-any.whl"
    python3 -m pip install --no-index --find-links "wheel" "/home/mks/python_package/wheel-0.42.0-py3-none-any.whl"
fi

# 安装依赖包
install_python_package_if_not_exists "qrcode" "qrcode-7.4.2-py3-none-any.whl"
install_python_package_if_not_exists "cryptography" "cryptography-41.0.7-cp37-abi3-manylinux_2_17_aarch64.manylinux2014_aarch64.whl"
install_python_package_if_not_exists "toml" "toml-0.10.2-py2.py3-none-any.whl"

# 删除本地包目录
rm -rf /home/mks/python_package
################本地更新python库结束

# 移动文件
# 修改主频 提升开机速度
cp -fp /root/etc/rk3328-roc-cc.dtb /boot/dtb/rockchip/rk3328-roc-cc.dtb
cp -fp /root/etc/rk3328-roc-cc.dts /boot/dtb/rockchip/rk3328-roc-cc.dts
cp -fp /root/etc/60-persistent-v4l.rules /lib/udev/rules.d/60-persistent-v4l.rules

#修改文件权限 (Q1Libre: fixed permissions - replaced broad chmod with targeted values)
find /home/mks/klipper_config -type d -exec chmod 755 {} \;
find /home/mks/klipper_config -type f -exec chmod 644 {} \;
chmod 700 /root/xindi/build/xindi
chown -R mks:mks /home/mks/klipper/
chown -R mks:mks /home/mks/klipper_config/

chmod 700 /usr/bin/makerbase-automount
# Q1Libre: set correct sudoers permissions
chmod 0440 /etc/sudoers.d/q1libre 2>/dev/null || true

# Q1Libre: install version marker
if [ -f /root/q1libre_version.txt ]; then
    echo "Q1Libre firmware installed" >> /var/log/q1libre_install.log
    echo "  Version: $(cat /root/q1libre_version.txt)" >> /var/log/q1libre_install.log
    echo "  Date: $(date)" >> /var/log/q1libre_install.log
fi

# Q1Libre Phase 2B: Upgrade Klipper to managed fork (danthi123/klipper q1-pro)
echo "Q1Libre: upgrading Klipper..."
# Stop klipper before replacing files
systemctl stop klipper.service 2>/dev/null || true
echo "Q1Libre: klipper stopped for upgrade"

# Back up existing klipper for rollback
if [ -d /home/mks/klipper ]; then
    rm -rf /home/mks/klipper.bak
    cp -a /home/mks/klipper /home/mks/klipper.bak \
        || echo "Q1Libre: WARNING - klipper backup failed, rollback not available"
    echo "Q1Libre: klipper backed up to /home/mks/klipper.bak"
fi

# New klipper tree is placed by dpkg from overlay — set ownership
chown -R mks:mks /home/mks/klipper

# Rebuild klippy-env with Python 3 for upstream Klipper v0.13
echo "Q1Libre: rebuilding klippy-env with Python 3..."
if [ -d /home/mks/klippy-env ]; then
    cp -a /home/mks/klippy-env /home/mks/klippy-env.bak
    rm -rf /home/mks/klippy-env
fi
/usr/bin/virtualenv -p python3 /home/mks/klippy-env 2>/dev/null
if [ -x /home/mks/klippy-env/bin/pip ]; then
    # Install from bundled wheels first (offline), fall back to PyPI if needed
    if [ -d /root/klippy-wheels ] && ls /root/klippy-wheels/*.whl >/dev/null 2>&1; then
        echo "Q1Libre: installing klippy packages from bundled wheels..."
        # Step 1: Bootstrap pip and setuptools FIRST
        # Stock virtualenv 15.1 bundles pip 18.1 which can't handle manylinux_2_17 wheels.
        # We must upgrade pip+setuptools before installing any native wheels.
        if ls /root/klippy-wheels/pip-*.whl >/dev/null 2>&1; then
            echo "Q1Libre: bootstrapping pip..."
            /home/mks/klippy-env/bin/pip install --no-index /root/klippy-wheels/pip-*.whl 2>&1 || true
        fi
        if ls /root/klippy-wheels/setuptools-*.whl >/dev/null 2>&1; then
            echo "Q1Libre: bootstrapping setuptools..."
            /home/mks/klippy-env/bin/pip install --no-index /root/klippy-wheels/setuptools-*.whl 2>&1 || true
        fi
        # Step 2: Install all remaining wheels (greenlet, cffi, pycparser,
        # Jinja2, markupsafe, pyserial, wrapt, aenum, numpy)
        echo "Q1Libre: installing klippy wheel packages..."
        /home/mks/klippy-env/bin/pip install --no-index --find-links /root/klippy-wheels \
            /root/klippy-wheels/greenlet-*.whl \
            /root/klippy-wheels/cffi-*.whl \
            /root/klippy-wheels/pycparser-*.whl \
            /root/klippy-wheels/Jinja2-*.whl \
            /root/klippy-wheels/MarkupSafe-*.whl \
            /root/klippy-wheels/pyserial-*.whl \
            /root/klippy-wheels/wrapt-*.whl \
            /root/klippy-wheels/aenum-*.whl \
            /root/klippy-wheels/numpy-*.whl \
            /root/klippy-wheels/scipy-*.whl \
            2>&1 || echo "Q1Libre: WARNING - some wheel installs had errors"
        # Step 3: Install python-can from source tarball (needs setuptools)
        if ls /root/klippy-wheels/python-can-*.tar.gz >/dev/null 2>&1; then
            echo "Q1Libre: installing python-can from sdist..."
            /home/mks/klippy-env/bin/pip install --no-index --find-links /root/klippy-wheels \
                /root/klippy-wheels/python-can-*.tar.gz 2>&1 || true
        fi
        rm -rf /root/klippy-wheels 2>/dev/null || true
        echo "Q1Libre: klippy packages installed"
    else
        echo "Q1Libre: installing klippy packages from PyPI..."
        /home/mks/klippy-env/bin/pip install -r /home/mks/klipper/scripts/klippy-requirements.txt \
            --quiet 2>/dev/null \
            || echo "Q1Libre: WARNING - klippy pip install had errors"
    fi
    chown -R mks:mks /home/mks/klippy-env
    echo "Q1Libre: klippy-env rebuilt with Python 3"
else
    echo "Q1Libre: WARNING - klippy-env rebuild failed, restoring backup"
    rm -rf /home/mks/klippy-env
    cp -a /home/mks/klippy-env.bak /home/mks/klippy-env
fi

# Reset git repo to the vendored upstream commit so update_manager sees a clean repo
KLIPPER_VERSION="v0.13.0-q1pro"
KLIPPER_SHA="a28e64cb3cfe62ba3ffdae7a3e3f4741bbaebbfe"
cd /home/mks/klipper
if [ -d .git ]; then
    git config user.email "q1libre@local"
    git config user.name "Q1Libre"
    git fetch origin 2>/dev/null || true
    if git cat-file -e $KLIPPER_SHA 2>/dev/null; then
        git reset --hard $KLIPPER_SHA 2>/dev/null
        echo "Q1Libre: klipper git reset to $KLIPPER_SHA"
    else
        echo "Q1Libre: target SHA not available (offline install), recreating repo..."
        rm -rf .git
    fi
fi
# If no .git exists (fresh install or offline fallback), create a clean repo
if [ ! -d .git ]; then
    git init 2>/dev/null
    git config user.email "q1libre@local"
    git config user.name "Q1Libre"
    git remote add origin https://github.com/danthi123/klipper.git 2>/dev/null || true
    git add -A 2>/dev/null
    git commit -m "Q1Libre: vendored klipper $KLIPPER_VERSION" --quiet 2>/dev/null
    git tag "$KLIPPER_VERSION" 2>/dev/null
    echo "Q1Libre: klipper repo initialized at $KLIPPER_VERSION"
fi
git remote set-url origin https://github.com/danthi123/klipper.git 2>/dev/null || true
git branch --set-upstream-to=origin/q1-pro master 2>/dev/null || true
chown -R mks:mks /home/mks/klipper/.git
cd /

# Pre-compile c_helper.so (takes ~30s on ARM, avoids timeout on first startup)
echo "Q1Libre: compiling klipper C helper..."
rm -f /home/mks/klipper/klippy/chelper/c_helper.so
cd /home/mks/klipper/klippy/chelper
/home/mks/klippy-env/bin/python -c "
import sys; sys.path.insert(0, '..')
import chelper; chelper.get_ffi()
" 2>&1 && echo "Q1Libre: c_helper compiled" || echo "Q1Libre: WARNING - c_helper compilation failed, will compile on first startup"
chown mks:mks /home/mks/klipper/klippy/chelper/c_helper.so 2>/dev/null || true
cd /

# Link timelapse macros config if available
if [ -f /home/mks/moonraker-timelapse/klipper_macro/timelapse.cfg ]; then
    ln -sf /home/mks/moonraker-timelapse/klipper_macro/timelapse.cfg /home/mks/klipper_config/timelapse.cfg
    chown -h mks:mks /home/mks/klipper_config/timelapse.cfg
fi

# Config migration not needed - fork supports all Qidi config options natively

echo "Q1Libre: Klipper upgrade complete"

# Q1Libre Phase 2A: Upgrade Moonraker
echo "Q1Libre: upgrading Moonraker..."
# Stop moonraker before replacing its files
systemctl stop moonraker.service 2>/dev/null || true
echo "Q1Libre: moonraker stopped for upgrade"

# Back up existing moonraker for rollback
if [ -d /home/mks/moonraker ]; then
    rm -rf /home/mks/moonraker.bak
    cp -a /home/mks/moonraker /home/mks/moonraker.bak \
        || echo "Q1Libre: WARNING - moonraker backup failed, rollback not available"
    echo "Q1Libre: moonraker backed up to /home/mks/moonraker.bak"
fi

# New moonraker tree is placed by dpkg from overlay — set ownership
chown -R mks:mks /home/mks/moonraker

# Reset git repo to the vendored upstream commit so update_manager sees a clean repo
# This allows the update manager to track and update Moonraker from the UI
MOONRAKER_VERSION="v0.10.0-19"
VENDORED_SHA="1ed102edfb34906115bfeb02712d4c839b93a7a9"
cd /home/mks/moonraker
if [ -d .git ]; then
    git config user.email "q1libre@local"
    git config user.name "Q1Libre"
    git fetch origin 2>/dev/null || true
    if git cat-file -e $VENDORED_SHA 2>/dev/null; then
        git reset --hard $VENDORED_SHA 2>/dev/null
        echo "Q1Libre: moonraker git reset to $VENDORED_SHA"
    else
        echo "Q1Libre: target SHA not available (offline install), recreating repo..."
        rm -rf .git
    fi
fi
# If no .git exists (fresh install or offline fallback), create a clean repo
if [ ! -d .git ]; then
    git init 2>/dev/null
    git config user.email "q1libre@local"
    git config user.name "Q1Libre"
    git remote add origin https://github.com/Arksine/moonraker.git 2>/dev/null || true
    git add -A 2>/dev/null
    git commit -m "Q1Libre: vendored moonraker $MOONRAKER_VERSION" --quiet 2>/dev/null
    git tag "$MOONRAKER_VERSION" 2>/dev/null
    echo "Q1Libre: moonraker repo initialized at $MOONRAKER_VERSION"
fi
# Exclude thirdparty symlink from git status
mkdir -p .git/info
echo "/thirdparty" >> .git/info/exclude
chown -R mks:mks /home/mks/moonraker/.git
cd /

# Create thirdparty symlink for module resolution (after git reset which removes it)
ln -sf /home/mks/moonraker/moonraker/thirdparty /home/mks/moonraker/thirdparty
chown -h mks:mks /home/mks/moonraker/thirdparty

# Link timelapse plugin into moonraker components
if [ -f /home/mks/moonraker-timelapse/component/timelapse.py ]; then
    ln -sf /home/mks/moonraker-timelapse/component/timelapse.py /home/mks/moonraker/moonraker/components/timelapse.py
    chown -h mks:mks /home/mks/moonraker/moonraker/components/timelapse.py
    echo "/moonraker/components/timelapse.py" >> /home/mks/moonraker/.git/info/exclude 2>/dev/null || true
fi

echo "Q1Libre: Moonraker upgrade complete"

# Q1Libre: Create printer_data symlinks for new Moonraker directory layout
mkdir -p /home/mks/printer_data
# Only create symlinks if they don't already exist as symlinks
[ -L /home/mks/printer_data/config ] || (rm -rf /home/mks/printer_data/config && ln -sf /home/mks/klipper_config /home/mks/printer_data/config)
[ -L /home/mks/printer_data/gcodes ] || (rm -rf /home/mks/printer_data/gcodes && ln -sf /home/mks/gcode_files /home/mks/printer_data/gcodes)
[ -L /home/mks/printer_data/logs ] || (rm -rf /home/mks/printer_data/logs && ln -sf /home/mks/klipper_logs /home/mks/printer_data/logs)
chown -h mks:mks /home/mks/printer_data /home/mks/printer_data/config /home/mks/printer_data/gcodes /home/mks/printer_data/logs

# Q1Libre: Enable mDNS/Avahi if available
if command -v avahi-daemon >/dev/null 2>&1; then
    systemctl enable avahi-daemon 2>/dev/null || true
    systemctl restart avahi-daemon 2>/dev/null || true
fi

#重启各类服务令
if [ "$method_value" = "0" ]; then
    echo "LAN only, only reloading and restarting klipper and moonraker"
    systemctl daemon-reload
    systemctl restart klipper.service
    systemctl restart moonraker.service
    sleep 12
    curl -X POST http://127.0.0.1:7125/printer/firmware_restart
else
    echo "Internet mode, performing full service operations"
    systemctl daemon-reload
    systemctl restart klipper.service
    systemctl restart moonraker.service
    systemctl enable frpc.service
    systemctl start frpc.service
    sleep 12
    curl -X POST http://127.0.0.1:7125/printer/firmware_restart
fi

sync
